# video_relay/robot/hikvision-ptz-agent@.service
[Unit]
Description=Hikvision ONVIF PTZ agent (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Develop/github_ws/hongshi_patrol_robot/video_relay/robot
ExecStart=/usr/bin/python3 /home/%i/Develop/github_ws/hongshi_patrol_robot/video_relay/robot/ptz_agent.py
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
